# systemd/tentacle-backlog.service
[Unit]
Description=Tentacle backlog review for OctopusGarden
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ryan
# Edit WorkingDirectory to match your clone path; %h expands to the user's home directory
EnvironmentFile=-%h/.config/tentacle/env
ExecStart=/usr/bin/python3 -m tentacle review-backlog
WorkingDirectory=%h/src/foundatron/tentacle
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-backlog

[Install]
WantedBy=multi-user.target
