# systemd/tentacle.timer
[Unit]
Description=Run tentacle research scout daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true
Unit=tentacle.service

[Install]
WantedBy=timers.target
